# flis_django/fabfile/supervisord-countries.conf
[unix_http_server]
file = %(countries_sandbox)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile = %(countries_sandbox)s/supervisor.log
logfile_maxbytes = 50MB
pidfile = %(countries_sandbox)s/supervisor.pid
user = edw
directory = %(countries_sandbox)s

[supervisorctl]
serverurl = unix://%(countries_sandbox)s/supervisor.sock

[program:flis_pt]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(pt-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8001
directory = %(countries_repo)s/pt
redirect_stderr = true
stdout_logfile = %(pt-instance_var)s/flis.log
startsecs = 4

[program:flis_cz]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(cz-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8002
directory = %(countries_repo)s/cz
redirect_stderr = true
stdout_logfile = %(cz-instance_var)s/flis.log
startsecs = 4

[program:flis_be]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(be-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8003
directory = %(countries_repo)s/be
redirect_stderr = true
stdout_logfile = %(be-instance_var)s/flis.log
startsecs = 4

[program:flis_ro]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(ro-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8004
directory = %(countries_repo)s/ro
redirect_stderr = true
stdout_logfile = %(ro-instance_var)s/flis.log
startsecs = 4

[program:flis_at]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(at-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8005
directory = %(countries_repo)s/at
redirect_stderr = true
stdout_logfile = %(at-instance_var)s/flis.log
startsecs = 4

[program:flis_ir]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(ir-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8006
directory = %(countries_repo)s/ir
redirect_stderr = true
stdout_logfile = %(ir-instance_var)s/flis.log
startsecs = 4

[program:flis_sk]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(sk-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8007
directory = %(countries_repo)s/sk
redirect_stderr = true
stdout_logfile = %(sk-instance_var)s/flis.log
startsecs = 4

[program:flis_hu]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(hu-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8008
directory = %(countries_repo)s/hu
redirect_stderr = true
stdout_logfile = %(hu-instance_var)s/flis.log
startsecs = 4

[program:flis_me]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(me-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8009
directory = %(countries_repo)s/me
redirect_stderr = true
stdout_logfile = %(me-instance_var)s/flis.log
startsecs = 4

[program:flis_mk]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(mk-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8010
directory = %(countries_repo)s/mk
redirect_stderr = true
stdout_logfile = %(mk-instance_var)s/flis.log
startsecs = 4

[program:flis_al]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(al-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8011
directory = %(countries_repo)s/al
redirect_stderr = true
stdout_logfile = %(al-instance_var)s/flis.log
startsecs = 4

[program:flis_xk]
command = %(countries_sandbox)s/bin/gunicorn_django -w 1 --user=edw --group=edw --log-level=debug  --log-file=%(xk-instance_var)s/supervisor_flis.log --bin 0.0.0.0:8012
directory = %(countries_repo)s/xk
redirect_stderr = true
stdout_logfile = %(xk-instance_var)s/flis.log
startsecs = 4

